# system_files/usr/lib/systemd/user/defenestra-nix-xdg-sync.service
[Unit]
# SPDX-License-Identifier: GPL-3.0-or-later
Description=Mirror Nix user profile XDG entries into ~/.local/share

[Service]
Type=oneshot
ExecStart=/usr/libexec/defenestra-nix-xdg-sync --user
Nice=10
